# system/uwsgi.service
[Unit]
Description=uWSGI Emperor
After=syslog.target

[Service]
Environment=statedir=/var/run/uwsgi
ExecStartPre=/bin/mkdir -p ${statedir}
ExecStart=/usr/bin/uwsgi --ini /etc/uwsgi/emperor.ini
Restart=always
KillSignal=SIGQUIT
Type=notify
StandardError=syslog
NotifyAccess=all

[Install]
WantedBy=multi-user.target